# init/s6/services.d/cync2mqtt/run
#!/usr/bin/with-contenv bashio

# =============================================================================
# ring-mqtt run script for s6-init               #
#
# This script automatically detects if it is running as the Home Assistant 
# addon or a standard docker environment and sets configuration variables as 
# appropriate for the detected environment.
# ==============================================================================

# Delay to keep logs messages from overlapping with s6 logs
sleep .5

if [ -f /data/options.json ]; then
    export RUNMODE=addon
else
    export RUNMODE=docker
fi

echo "Running cync2mqtt..."
exec /cync2mqtt/bin/cync2mqtt /config/cync2mqtt/cync_mesh.yaml
